# deploy/systemd/momentum-premarket-eval.service
[Unit]
Description=Momentum premarket (pm1) shadow track rebuild + re-validate

[Service]
Type=oneshot
WorkingDirectory=/mnt/c/code/momentum
# DATABASE_URL loads from the repo's gitignored .env via research.premarket_eval (python-dotenv)
ExecStart=/home/philip/.venvs/momentum/bin/python -m research.premarket_eval
StandardOutput=append:/mnt/c/code/momentum/data/premarket_eval.log
StandardError=append:/mnt/c/code/momentum/data/premarket_eval.log
